# ops/launchd/com.openread.dev.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.openread.dev.cleanup

  Runs /cleanup --scheduled=weekly every Sunday at 03:00 local time.

  Install:
    cp ops/launchd/com.openread.dev.cleanup.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.openread.dev.cleanup.plist

  Unload:
    launchctl unload ~/Library/LaunchAgents/com.openread.dev.cleanup.plist

  Reference: docs/probe-driven-workflow-implementation.md §7.5
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openread.dev.cleanup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$OPENREAD_REPO_ROOT" &amp;&amp; /usr/bin/env pnpm exec tsx scripts/skills/cleanup.ts --scheduled=weekly --json</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Override via `launchctl setenv OPENREAD_REPO_ROOT <path>` before load. -->
        <key>OPENREAD_REPO_ROOT</key>
        <string>${HOME}/ClaudCode/openread.ai.probe-infra</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Sunday 03:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Defer (don't pile up) if Mac was asleep at the scheduled time. -->
    <key>StartCalendarIntervalRunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>${HOME}/.openread-dev/cleanup.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.openread-dev/cleanup.log</string>

    <!-- Don't re-spawn on exit; this is a scheduled one-shot. -->
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <false/>

    <!-- Nice priority: be a good citizen while the developer is working. -->
    <key>Nice</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Allow up to 30 minutes of wall-clock per run before launchd SIGTERMs. -->
    <key>ExitTimeOut</key>
    <integer>1800</integer>
</dict>
</plist>
